# templates/etc/profile.d/java-openjdk.sh
JAVA_HOME="{{ java_openjdk_install_folder }}/java"
PATH="$JAVA_HOME/bin:$PATH"
